# dot_config/systemd/user/hyprlock.service.tmpl
[Unit]
Description=Hyprlock Screen Lock
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/hyprlock
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
