# dot_config/systemd/user/hyprland.service.tmpl
[Unit]
Description=Hyprland Wayland Compositor
After=graphical-session.target

[Service]
ExecStart=/usr/bin/Hyprland
Restart=always

[Install]
WantedBy=default.target
